# cookbooks/node-upstart/templates/default/nodeSites.conf.erb
#!upstart
# nodeSites.conf
# Generated by Chef for <%= node['hostname'] %>.
# Local modifications will be overwritten.
# This script is to be placed in /etc/init to work with upstart.
#
# http://clock.co.uk/tech-blogs/upstart-and-nodejs

description 'Start all node instances at boot'
author      "Vincent Mac <vincent@simplicity.io>"
version     "1.0.1"

start on (local-filesystems and net-device-up)

task

script
  for file in `ls /etc/node/*.conf` ; do
    filename=`basename ${file%.*}`
    start node NAME=$filename
  done
end script