# cookbooks/node-upstart/templates/default/node.conf.erb
#!upstart
# node.conf
# Generated by Chef for <%= node['hostname'] %>.
# Local modifications will be overwritten.
# This script is to be placed in /etc/init to work with upstart.
#
# Internally the 'initctl' command is used to manage:
# initctl help

# http://clock.co.uk/tech-blogs/upstart-and-nodejs

description "node upstart script"
author      "Vincent Mac <vincent@simplicity.io>"
version     "1.0.1"

start on (local-filesystems and net-device-up)
stop on shutdown

instance "Node - $NAME"

respawn
respawn limit 5 60

script
  . /etc/node/$NAME.conf
  exec sudo -u node NODE_ENV=${NODE_ENV} /usr/bin/n use ${NODE_VERSION} ${NODE_PATH} >> ${LOG_PATH} 2>&1
end script